# build/ferretdb/evaluation/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:postgres]
command=/usr/local/bin/docker-entrypoint.sh postgres
; fast shutdown https://www.postgresql.org/docs/current/server-shutdown.html
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:ferretdb]
command=/usr/local/bin/ferretdb
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
